# VUE2/Nullsoft-Installer/jre/lib/task.xml
<?xml version="1.0" ?>
<!--
This sample schedules a task to start on a daily basis.
-->
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
    <RegistrationInfo>
        <Date>2009-01-01T00:00:00</Date>
        <Author>Java</Author>
        <Description>In OEMUpdate, jusched starts when users login</Description>
    </RegistrationInfo>
    <Triggers>
        <LogonTrigger>
            <StartBoundary>2005-01-01T00:00:00</StartBoundary>
            <EndBoundary>2015-01-01T00:00:00</EndBoundary>
            <Enabled>true</Enabled>
            <Delay>PT120S</Delay>
        </LogonTrigger>
    </Triggers>
    <Principals>
        <Principal>
            <GroupId>Users</GroupId>
        </Principal>
    </Principals>
    <Settings>
        <Enabled>true</Enabled>
        <AllowStartOnDemand>true</AllowStartOnDemand>
        <AllowHardTerminate>true</AllowHardTerminate>
        <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
        <DeleteExpiredTaskAfter>PT0S</DeleteExpiredTaskAfter>
        <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    </Settings>
    <Actions>
        <Exec>
            <Command>%WINDIR%\System32\jusched.exe</Command>
        </Exec>
    </Actions>
</Task>

